# gui/server/clean.crontab
0 * * * * $ATTRACTGUI/server/clean.cron >& /dev/null
