# gui/server/clean-tum.crontab
0 * * * * $ATTRACTGUI/server/clean-tum.cron >& /dev/null
